# debian/postinst
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload || true
    systemctl enable nvidia-driver-manager.timer >/dev/null 2>&1 || true
fi

exit 0
